# crontab
# m   h   dom mon dow command
  *   *   *   *   *   /pf_update.sh
